# caddy/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile
autostart=true
autorestart=true